# root/etc/services.d/bilbili-helper/run
#!/usr/bin/with-contenv bash

# 启动定时任务,默认一天执行一次
# 启动bilbili-helper
java -jar /app/BILIBILI-HELPER.jar ${DEDEUSERID} ${SESSDATA} ${BILI_JCT}

echo -e "$(date +"%m/%d %H:%M:%S")"
echo ==============今天已签到完毕，下次签到时间为${TASK}后==============
sleep $TASK
